# overlay/usr/lib/systemd/system/first_time_setup.service
[Unit]
Description=First Time Setup

[Service]
ExecStart=/usr/local/sbin/first_time_setup.sh

[Install]
WantedBy=multi-user.target
