# overlay/usr/lib/systemd/system/resize_rootfs.service
[Unit]
Description=Resize rootfs

[Service]
ExecStart=/usr/local/sbin/resize_rootfs.sh

[Install]
WantedBy=multi-user.target
