# overlay/usr/lib/systemd/system/setup_audio.service
[Unit]
Description=Setup audio

[Service]
ExecStart=/usr/local/sbin/setup_audio.sh

[Install]
WantedBy=multi-user.target
