# overlay/usr/lib/systemd/system/usb_tethering.service
[Unit]
Description=USB tethering
Before=dhcpd4.service

[Service]
ExecStart=/usr/local/sbin/usb_tethering.sh

[Install]
WantedBy=multi-user.target
